# packaging/controlplane/gateway-agent/gateway-agent.service
# vi:syntax=systemd
[Unit]
Description=naisdevice gateway agent

[Service]
Restart=always
EnvironmentFile=/etc/default/gateway-agent
StandardOutput=append:/var/log/naisdevice/gateway-agent.json
StandardError=append:/var/log/naisdevice/gateway-agent.json
ExecStart=/usr/bin/gateway-agent

[Install]
WantedBy=multi-user.target

