# packaging/controlplane/prometheus-agent/prometheus-agent.service
# vi:syntax=systemd
[Unit]
Description=naisdevice prometheus agent

[Service]
Restart=always
RestartSteps=4
RestartMaxDelaySec=30s
EnvironmentFile=/etc/default/prometheus-agent
ExecStart=/usr/bin/prometheus-agent

[Install]
WantedBy=multi-user.target
