# ansible/roles/apiserver/files/naisdevice-health-checker.timer
[Unit]
Description=Naisdevice health checker every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
