# packaging/controlplane/apiserver/apiserver.service
# vi:syntax=systemd
[Unit]
Description=naisdevice apiserver

[Service]
Restart=always
EnvironmentFile=/etc/default/apiserver
ExecStart=/usr/bin/apiserver

[Install]
WantedBy=multi-user.target
